# deploy/hospital-web.service
[Unit]
Description=Hospital CMS Web
After=network.target hospital-api.service
Wants=network.target

[Service]
Type=simple
User=hospitalcms
Group=hospitalcms
WorkingDirectory=/opt/hospital-cms
EnvironmentFile=/etc/hospital-cms/hospital.env
ExecStart=/usr/bin/env pnpm --dir /opt/hospital-cms/apps/web start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
